# data/systemd/timelapse.service.in
[Unit]
Description=Image Capture Microservice
AssertPathExists=@sysconfdir@

[Service]
#User=pi
#Group=pi
ExecStart=@bindir@/timelapse --config @sysconfdir@/config.ini

[Install]
WantedBy=default.target
